# stormymoose.service
[Unit]
Description=stormymoose jenkins configuration

[Service]
ExecStart=/usr/local/sbin/stormymoose

[Install]
WantedBy=multi-user.target
